# more.service
[Unit]
Description=MORE - Monitoreo del sistema

[Service]
Type=oneshot
ExecStart=/opt/more/datos.sh
